# roles/minecraft/templates/minecraft.socket.j2
[Socket]
ListenStream = 25565
Accept=no
Service=minecraft.service

[Install]
WantedBy = sockets.target